# etc/systemd/system/mycustom.service
[Unit]
Description=Startup

[Service]
ExecStart=/usr/bin/iw dev wlp3s0 set power_save off

[Install]
WantedBy=multi-user.target
